# templates/profile.d/karaf-profile-deb.erb
# File managed by Puppet
if [ -d <%= scope.lookupvar('karaf::karaf_dir') %>/bin ]; then
	PATH=$PATH:<%= scope.lookupvar('karaf::karaf_dir') %>/bin
fi
